# Silk/silk.service
[Unit]
Description="Silk API service"

[Service]
User=root
WorkingDirectory=/home/tissue/Tissue/Sink
ExecStart=/home/tissue/.cargo/bin/silk

[Install]
WantedBy=multi-user.target
