# Silk/silk.service
[Unit]
Description="Silk API service"

[Service]
User=tissue
WorkingDirectory=/home/tissue/Tissue/Silk
ExecStart=/home/tissue/.cargo/bin/silk

[Install]
WantedBy=multi-user.target
