# Silk/silk.service
[Unit]
Description="Silk API service"

[Service]
User=tissue
WorkingDirectory=/home/tissue/Tissue/Sink
ExecStart=cargo run --release

[Install]
WantedBy=multi-user.target
